# infra/systemd/bullet-in-fail-notify.service
[Unit]
Description=bullet-in cycle failure alert (Discord)

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/bullet-in
EnvironmentFile=/home/ubuntu/bullet-in/.env
ExecStart=/home/ubuntu/.local/bin/uv run python -c "from bullet_in.notify import send_alert; send_alert('bullet-in 회차 실패 (systemd)', 'bullet-in.service 실패 — VM 에서 journalctl -u bullet-in -n 100 확인', color=0xE74C3C)"
